# dot_config/systemd/user/wallpaper.timer.tmpl
[Unit]
Description=Timer for wallpaper rotation
PartOf=graphical-session.target

[Timer]
OnUnitActiveSec={{ .systemd.user.wallpaper.timer.period }}
